# ind-efd/root/etc/systemd/system/sepl-modem.service
[Unit]
Description=Check Internet Connection
After=network.target auditd.service

[Service]
ExecStart=/opt/sbin/check-internet-connection.sh
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target

