# ind-efd/root/etc/cron.d/efd_ping
# Perform an EFD Ping everying minute
* * * * * root /opt/sbin/efd-ping.sh
