# ind-efd/root/etc/cron.d/efd_chrony_watchdog
#! Periodically call chrony watchdog
#! every 10 minutes
*/10 * * * * root /opt/sbin/chrony_watchdog.sh
