# ind-efd/root/etc/cron.d/ind-efd
##
## cron.d/ind-efd -- Perform all IND-EFD specific cron jobs
##

## By default, run at 00:57 on every Sunday, but do nothing unless the day of
## the month is less than or equal to 7. Thus, only run on the first Sunday of
## each month. crontab(5) sucks, unfortunately, in this regard; therefore this
## hack (see #380425).
#57 0 * * 0 root if [ -x /usr/share/mdadm/checkarray ] && [ $(date +\%d) -le 7 ]; then /usr/share/mdadm/checkarray --cron --all --idle --quiet; fi

# do an EFD Ping everying minute
* * * * * root /opt/sbin/efd-ping.sh

